# roles/b2bua/templates/wazo-sysconfd.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/xivo-sysconfd --listen-addr {{ b2bua_listen_address }} --listen-port {{ b2bua_port_sysconfd }}
