# roles/b2bua/templates/xivo-sysconfd.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/xivo-sysconfd --listen-addr 0.0.0.0 --listen-port {{ b2bua_port_sysconfd }}